# deploy/vllm-hymt2.service
# Hy-MT2-7B — быстрый переводной контур (виджет/фрагмент), GPU1, порт 8005.
# Преемник HY-MT1.5-7B (та же arch hunyuan_v1_dense). Принят 2026-06-19 по COMET-A/B
# на корпусе заказчика (COMETKiwi-22, 300 сегм.): средний 0.7790 vs 0.7716 у HY-MT1.5
# и 0.7516 у прод-Qwen3.5; p10 0.618 (крепче хвост). bf16 — НЕ FP8: FP8-чекпойнт на
# A100 (Ampere, Marlin-эмуляция) выдаёт мусорные логиты. Память 0.5 — контуру коротких
# выделений большой KV-пул не нужен (раньше HY-MT1.5 держал 0.85 ≈ 34 ГБ зря).
[Unit]
Description=vLLM Hy-MT2-7B (GPU1, :8005)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=CUDA_VISIBLE_DEVICES=1
Environment=NCCL_P2P_DISABLE=1
Environment=NCCL_IB_DISABLE=1
Environment=CUDA_HOME=/usr/local/cuda
Environment=PATH=/root/services/vllm-qwen32b/.venv/bin:/usr/local/cuda/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# GPU1 — NUMA node0 ({1,2} ↔ node0)
ExecStart=/usr/bin/numactl --cpunodebind=0 --membind=0 \
    /root/services/vllm-qwen32b/.venv/bin/vllm serve /root/models/Hy-MT2-7B \
    --served-model-name hy-mt2-7b \
    --host 127.0.0.1 \
    --port 8005 \
    --gpu-memory-utilization 0.5 \
    --max-model-len 8192 \
    --max-num-seqs 64
Restart=always
RestartSec=10
TimeoutStartSec=900
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
