# deploy/vllm-visual-embedding.service
# Визуальный эмбеддер Qwen3-VL-Embedding-8B — GPU2, порт 8007 (§ 12.1 шаг 4).
# Визуальный контур (эмбеддер + реранкер) на свободной GPU2; текстовый — на GPU4.
[Unit]
Description=vLLM Qwen3-VL-Embedding-8B (GPU2, :8007)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=CUDA_VISIBLE_DEVICES=2
Environment=NCCL_P2P_DISABLE=1
Environment=NCCL_IB_DISABLE=1
Environment=CUDA_HOME=/usr/local/cuda
# Отдельный venv с последним vLLM: 0.22.1 давал NaN в мультимодальном
# pooling-пути этой модели (строковый /v1/embeddings работал)
Environment=PATH=/root/services/vllm-visual/.venv/bin:/usr/local/cuda/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# GPU2 — NUMA node0
ExecStart=/usr/bin/numactl --cpunodebind=0 --membind=0 \
    /root/services/vllm-visual/.venv/bin/vllm serve /root/models/Qwen3-VL-Embedding-8B \
    --served-model-name qwen3-vl-embedding-8b \
    --runner pooling \
    --host 127.0.0.1 \
    --port 8007 \
    --gpu-memory-utilization 0.6 \
    --max-model-len 16384 \
    --dtype float16 \
    --enforce-eager
# --enforce-eager ОБЯЗАТЕЛЕН: скомпилированные ядра vision-башни (torch.compile)
# в vLLM 0.22.1 на Ampere дают NaN на реальных страницах (зависело от содержимого
# пикселей); eager стабилен на полном разрешении. float16: bf16-путь NaN'ился
# даже на мини-картинках. Перепроверить оба флага при апгрейде vLLM.
Restart=always
RestartSec=10
TimeoutStartSec=900
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
